# bad_env-web-1.conf
start on starting bad_env-web
stop on stopping bad_env-web
respawn

env PORT=6000
env PORT='6000'

setuid me

chdir /home/vagrant/test_project

exec bundle exec puma -p $PORT -e $RACK_ENV
