# bad_env-web.conf
start on starting bad_env
stop on stopping bad_env
