# good_env-web-2.conf
start on starting good_env-web
stop on stopping good_env-web
respawn

env PORT='6000'
env PORT=6001

setuid me

chdir /home/vagrant/test_project

exec bundle exec puma -p $PORT -e $RACK_ENV
